# systemd/waid.service
[Unit]
Description=waid desktop activity tracker
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/waid run --config %h/.config/waid/config.yaml
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
