# src/waid/resources/systemd/waid.service.in
[Unit]
Description=waid desktop activity tracker
After=graphical-session.target

[Service]
Type=simple
ExecStart={{ exec_start }}
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
